# overlay/etc/local.d/00_root_permissions.start
#!/bin/bash

__self__="$(realpath "$0")"

chown -R 0:0 /root
chmod 700 /root /root/.ssh
chmod 600 /root/.ssh/authorized_keys

mv "$__self__" "$__self__.ran"
