# overlay/etc/local.d/10_create_folders.start
#!/bin/bash

__self__="$(realpath "$0")"

mkdir -p /var/lib/misc  # for dnsmasq

mv "$__self__" "$__self__.ran"
