# overlay/etc/init.d/05_root_permissions
#!/sbin/openrc-run

chown -R 0:0 /root
chmod 700 /root /root/.ssh
chmod 600 /root/.ssh/authorized_keys
rm -vf /etc/init.d/05_root_permissions
